# deploy/macbook/com.wca.research.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.wca.research</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/andrewdoherty/Desktop/Coding/World Cup Alpha/deploy/macbook/research_cycle.sh</string>
  </array>
  <key>StartInterval</key><integer>3600</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/andrewdoherty/Desktop/Coding/World Cup Alpha/logs/research-cycle.log</string>
  <key>StandardErrorPath</key><string>/Users/andrewdoherty/Desktop/Coding/World Cup Alpha/logs/research-cycle.log</string>
</dict></plist>
